# etc/systemd/system/home-address-watcher.path
# SPDX-License-Identifier: GPL-3.0-or-later

[Path]
PathModified=/etc/home-address-ipv6

[Install]
WantedBy=multi-user.target
